# packages/boards/nvidia/tegra/src/common/etc/systemd/system/var.mount
[Unit]
Description=Variable State Directory
Before=local-fs.target

[Mount]
What=tmpfs
Where=/var
Type=tmpfs
Options=mode=755

[Install]
WantedBy=local-fs.target
